# system/crontab
0 0 * * * sh /valhalla/system/3rdparty.sh
0 0 * * * service dnsmasq restart
